# supervisor.conf
[program:redis]
command=redis-server /etc/redis/redis.conf
